# cron/plantwatch.cron
# Timelapse capture every 5 minutes (actual capture controlled by app settings)
*/5 * * * * __INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/cam.py --timelapse >> __INSTALL_DIR__/cam_timelapse.log 2>&1

# Watering reminder check once daily at 16:45
45 16 * * * __INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/check_watering.py >> __INSTALL_DIR__/check_watering.log 2>&1
