# systemd/plantwatch-bot.service
[Unit]
Description=PlantWatch Telegram Bot Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=__INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/bot_listener.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
